# V3-DLP-Fan-Guide-Complete-20260428-161604/docker/collect-waits.cron
*/30 * * * * php /var/www/html/jobs/collect-waits.php >> /proc/1/fd/1 2>> /proc/1/fd/2
